# config/docker/services/depot_services.conf
[supervisord]
nodaemon=true
directory=/tmp
pidfile=/tmp/supervisord.pid

[program:depot]
command=/usr/local/python3/bin/uwsgi --ini /etc/uwsgi-depot.ini
autostart=true
autorestart=unexpected
stdout_logfile=/logs/depot_uwsgi.log
redirect_stderr=true
startretries=0
exitcodes=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
